# build/fbcpX.service
[Unit]
Description=Fast Framebuffer Copy Service for the Raspberry Pi
StartLimitIntervalSec=10

[Service]
Type=simple
EnvironmentFile=/etc/fbcpx/fbcpX.conf
Restart=always
RestartSec=1
StartLimitBurst=2
ExecStart=/usr/local/bin/fbcpx

[Install]
WantedBy=default.target
